# fleet-unit-files/cron/dockermi.timer
[Unit]
Description=Timer for dockermi service
Requires=docker.service
[Timer]
OnCalendar=daily
[X-Fleet]
MachineOf=dockermi.service
